# crontab
0 0 * * * cd /path/to/parole-hearing && ./run.sh
